# roles/cart/templates/cart.service.j2
[Unit]
Description = Cart Service
[Service]
User=roboshop
// highlight-start
Environment=REDIS_HOST="{{Redis_Host}}"
Environment=CATALOGUE_HOST="{{Catalogue_Host}}
Environment=CATALOGUE_PORT=80
// highlight-end
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target